# 1744833742-freedomofpress-securedrop/securedrop/debian/config/lib/systemd/system/securedrop-noble-migration-check.timer
[Unit]
Description=Check noble migration readiness

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=5m

[Install]
WantedBy=timers.target
