# 1744833742-freedomofpress-securedrop/securedrop/debian/config/lib/systemd/system/securedrop-cleanup-ossec.timer
[Unit]
Description=Cleanup OSSEC diff queue

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=5m

[Install]
WantedBy=timers.target
