# 1744833742-freedomofpress-securedrop/securedrop/debian/config/lib/systemd/system/securedrop-noble-migration-upgrade.timer
[Unit]
Description=Run noble migration

[Timer]
OnBootSec=3m
OnUnitInactiveSec=3m
Persistent=true

[Install]
WantedBy=timers.target
