# 1744833742-freedomofpress-securedrop/securedrop/debian/config/lib/systemd/system/securedrop-reboot-required.timer
[Unit]
Description=Touch reboot-required file every 12 hours

[Timer]
OnCalendar=*-*-* 0/12:00:00

[Install]
WantedBy=timers.target
