# 1744833742-freedomofpress-securedrop/securedrop/debian/config/lib/systemd/system/securedrop-noble-migration-check.service
[Unit]
Description=Check noble migration readiness

[Service]
Type=exec
ExecStart=/usr/bin/securedrop-noble-migration-check
User=root
